# auto-update/claude-desktop-update.timer
[Unit]
Description=Hourly check for new Claude Desktop builds

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
RandomizedDelaySec=5min
Persistent=true

[Install]
WantedBy=timers.target
